# insectius-monitor.service
[Unit]
Description=Insectius Monitor - Sensor Data Collection and API Sync
After=network.target bluetooth.target graphical.target
Wants=bluetooth.target
# Esperar a que Bluetooth esté completamente listo
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Environment="HOME=/root"
WorkingDirectory=/opt/insectius-monitor
# Esperar a que Bluetooth esté listo antes de iniciar
ExecStartPre=/bin/sleep 5
ExecStart=/opt/insectius-monitor/insectius-monitor
Restart=always
RestartSec=15

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=insectius-monitor

# Capabilities for Bluetooth access
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=graphical.target
